# ek3030e-modbus/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root

; ── Web-Dashboard (Node.js) ───────────────────────────────────────────────────
[program:webapp]
command=node /app/dist/index.js
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",PORT="%(ENV_PORT)s",MODBUS_BRIDGE_URL="%(ENV_MODBUS_BRIDGE_URL)s"
priority=10

; ── Python Modbus-Bridge ──────────────────────────────────────────────────────
[program:modbus_bridge]
command=python3 /app/server/modbus_bridge.py --host %(ENV_MODBUS_BRIDGE_HOST)s --port %(ENV_MODBUS_BRIDGE_PORT)s
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=5
